# docker/root/etc/services.d/gcsfuse/run
#!/usr/bin/with-contenv bash

exec gcsfuse --foreground -o allow_other $BUCKET_MEDIA /data
